# etc/profile.d/git-aliases
#!/bin/bash

alias g="git"
alias ga="git add"
alias gb="git branch -vv"
alias gb+="git branch -vv --all"
alias gc="git commit --gpg-sign"
alias gc+="git commit --gpg-sign --all"
alias gco="git checkout"
alias gco+="git checkout -b"
alias gd="git diff"
alias gd+="git diff --cached"
alias gl="git log --graph --oneline --decorate"
alias gl+="git log --graph --oneline --decorate --stat"
alias gf="git fetch --prune"
alias gp="git pull --prune"
alias gpu="git push"
alias gpu+="git push --force"
alias gpuo="git push --set-upstream origin"
alias gm="git merge --no-ff"
alias gs="git status --short"
alias gr="git reset HEAD"
alias gr+="git reset --hard HEAD"
alias gr++="git reset --hard HEAD~1"
alias git-prune='git fetch --prune && for branch in `git branch -vv | grep ": gone]" | awk "{print \\$1}"`; do git branch --delete --force $branch; done'
